# crontab
# Run every minute
*/1 * * * * bundle exec rake search:refresh